# ops/backup/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 3 * * * /opt/backup/backup_to_s3.sh >> /var/log/backup.log 2>&1
