# ops/deploy/systemd-family-ai.service
[Unit]
Description=Family AI Companion Docker Stack
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/family-ai
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
